# workspace/crontab/erp
* * * * * laradock /usr/bin/php /var/www/erp/artisan schedule:run --env=local >> /dev/null 2>&1
